# debian/pellmon.prerm
#!/bin/sh

set -e


if [ -x "/etc/init.d/pellmonweb" ]; then
	invoke-rc.d pellmonweb stop || true
fi

if [ -x "/etc/init.d/pellmonsrv" ]; then
	invoke-rc.d pellmonsrv stop || true
fi

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
