# debian/pellmon.postrm
#!/bin/sh

set -e



if [ "$1" = "purge" ] ; then
    deluser --quiet pellmonsrv dialout > /dev/null || true
    deluser --quiet --system pellmonsrv > /dev/null || true
    deluser --quiet --system pellmonweb > /dev/null || true
    delgroup --quiet --system pellmonsrv > /dev/null || true
    delgroup --quiet --system pellmonweb > /dev/null || true
    # Cleanup /var/run/pellmon
    f=/var/run/pellmon
    if [ -d "$f" ]; then
        rmdir "$f" || { chown root:root "$f" && chmod 00700 "$f"; }
    fi
fi

if [ "$1" = "purge" ] ; then
	update-rc.d pellmonsrv remove >/dev/null || true
	update-rc.d pellmonweb remove >/dev/null || true
fi

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

case "$1" in
  purge)
    rm -rf /var/log/pellmonsrv
    rm -rf /var/log/pellmonweb
    rm -rf /var/lib/pellmon
  ;;
esac

