# debian/pellmon.postinst
#!/bin/sh

set -e

case "$1" in
    configure)

# create the group if it isn't already there
	addgroup --quiet --system pellmonsrv
	addgroup --quiet --system pellmonweb

# ...create users if they aren't already there...
	adduser --system --group --no-create-home --home /nonexistent pellmonsrv \
	    --quiet --gecos "pellmonsrv pellet logger daemon"
	adduser --system --group --no-create-home --home /nonexistent pellmonweb \
	    --quiet --gecos "pellmonweb pellet webserver daemon"

    # Ask the bus to reload the config file
    if [ -x "/etc/init.d/dbus" ]; then
      invoke-rc.d dbus force-reload || true
    fi

# Add the pellmonsrv user to the dialout group, some plugins need it
    adduser --quiet pellmonsrv dialout
	;;
esac


# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

if [ -x "/etc/init.d/pellmonsrv" ]; then
	if [ ! -e "/etc/init/pellmonsrv.conf" ]; then
		update-rc.d pellmonsrv defaults >/dev/null
    		invoke-rc.d pellmonsrv start || true
	fi
fi

if [ -x "/etc/init.d/pellmonweb" ]; then
	if [ ! -e "/etc/init/pellmonweb.conf" ]; then
		update-rc.d pellmonweb defaults >/dev/null
	    	invoke-rc.d pellmonweb start || true
	fi
fi

exit 0


